# scripts/sniper.service
[Unit]
Description=Server for sniper

[Service]
ExecStart=%h/.local/bin/sniper
NonBlocking=True